# debian/postinst
#!/bin/bash
# postinst script for astroberry-os-sysmod

if [ "$1" = "configure" ]; then

    # Set required permissions if cmake failed
    if [ -e /etc/NetworkManager/system-connections/Hotspot.nmconnection ]; then chmod 600 /etc/NetworkManager/system-connections/Hotspot.nmconnection; fi

    # Set default desktop backgrounds
    update-alternatives --install /usr/share/images/desktop-base/desktop-background desktop-background /usr/share/astroberry-artwork/backgrounds/sci-fi-planet-darkness-night_1920x1080_bw.jpg 80
    update-alternatives --install /usr/share/images/desktop-base/login-background.svg desktop-login-background /usr/share/astroberry-artwork/backgrounds/milkyway-galaxy-center-and-its-companions_1920x1080.jpg 80
    update-alternatives --install /usr/share/images/desktop-base/desktop-grub.png desktop-grub /usr/share/astroberry-artwork/grub/milkyway-galaxy-center-and-its-companions_1920x1080.png 80

    # Enable services
    if [ ! ischroot ]; then
        if [ ! -e /etc/systemd/system/multi-user.target.wants/astroberry-indiwebmanager.service ]; then ln -fs /etc/systemd/system/astroberry-indiwebmanager.service /etc/systemd/system/multi-user.target.wants/; fi
        if [ ! -e /etc/systemd/system/multi-user.target.wants/astroberry-desktop-web.service ]; then ln -fs /etc/systemd/system/astroberry-desktop-web.service /etc/systemd/system/multi-user.target.wants/; fi
        if [ ! -e /etc/systemd/system/multi-user.target.wants/astroberry-desktop-vnc.service ]; then ln -fs /etc/systemd/system/astroberry-desktop-vnc.service /etc/systemd/system/multi-user.target.wants/; fi
        if [ ! -e /etc/systemd/system/multi-user.target.wants/astroberry-websockify.service ]; then ln -fs /etc/systemd/system/astroberry-websockify.service /etc/systemd/system/multi-user.target.wants/; fi
    else
        systemctl enable astroberry-indiwebmanager.service && systemctl start astroberry-indiwebmanager.service
        systemctl enable astroberry-desktop-web.service && systemctl start astroberry-desktop-web.service
        systemctl enable astroberry-desktop-vnc.service && systemctl start astroberry-desktop-vnc.service
        systemctl enable astroberry-websockify.service && systemctl start astroberry-websockify.service
    fi

    # Make sure that Astroberry repository is up to date
    if [ ! -e /etc/apt/sources.list.d/astroberry.sources ]; then
        curl -fsSL https://astroberry.io/debian/astroberry.asc | gpg --dearmor -o /etc/apt/keyrings/astroberry.gpg

        cat <<EOF > /etc/apt/sources.list.d/astroberry.sources
# Do not edit this file. It will be overwritten by each upgrade.
# Instead, add a new sources file in /etc/apt/sources.list.d/
Types: deb
URIs: https://astroberry.io/debian/
Suites: trixie
Architectures: arm64 amd64
Components: main restricted
Signed-By: /etc/apt/keyrings/astroberry.gpg
EOF
    fi

    # Pin Astroberry repository
    if [ ! -e /etc/apt/preferences.d/astroberry-pin ]; then
    cat <<EOF > /etc/apt/preferences.d/astroberry-pin
# Do not edit this file. It will be overwritten by each upgrade.
# Instead, add a new config file in /etc/apt/sources.list.d/
Package: *
Pin: origin astroberry.io
Pin-Priority: 900
EOF
    fi

fi
