# debian/postinst
#!/bin/bash
# postinst script for astroberry-os-sysmod

set -e

if [ "$1" == "configure" ]; then

	# Install astroberry-manager
	if [ ! -e /opt/astroberry-manager ]; then python -m venv /opt/astroberry-manager; fi
	/opt/astroberry-manager/bin/pip install --upgrade astroberry_manager@git+https://github.com/astroberry-official/astroberry-manager

	# Set wireless regulatory domain
	if [ -e /boot/firmware/cmdline.txt ] && [ -z "$(grep cfg80211.ieee80211_regdom /boot/firmware/cmdline.txt)" ]; then
		sed -i -e "s/\s*cfg80211.ieee80211_regdom=\S*//" -e "s/\(.*\)/\1 cfg80211.ieee80211_regdom=GB/" /boot/firmware/cmdline.txt
	fi

	# Disable graphical login
	# update-rc.d lightdm disable

	# Enable TigerVNC
	# ln -s /etc/systemd/system/vncserver.service /etc/systemd/system/multi-user.target.wants/

	# Enable INDI Web Manager
	# ln -s /etc/systemd/system/indiwebmanager.service /etc/systemd/system/multi-user.target.wants/

	# Enable Websockify
	# ln -s /etc/systemd/system/websockify.service /etc/systemd/system/multi-user.target.wants/

	# Enable Astroberry Manager
	# ln -s /etc/systemd/system/astroberry-manager.service /etc/systemd/system/multi-user.target.wants/

fi
