# debian/postinst
#!/bin/bash
# postinst script for astroberry-os-sysmod

set -e

if [ "$1" == "configure" ]; then

# Set cloud-init
if [ -e /boot/firmware/meta-data ]; then sed -i "s/instance_id: .*/instance_id: astroberryos-image/g" /boot/firmware/meta-data; fi

if [ -e /boot/firmware/user-data ]; then cp /boot/firmware/user-data /boot/firmware/user-data.default; fi
cat <<EOF > /boot/firmware/user-data
#cloud-config
hostname: astroberry
manage_etc_hosts: true

users:
- name: astroberry
  groups: users,adm,dialout,audio,netdev,video,plugdev,cdrom,games,input,gpio,spi,i2c,render,sudo
  shell: /bin/bash
  lock_passwd: false
  passwd: \$5\$MRvtBeRUfh\$QGv7UN3cnLPJ00paa82TjuwrjqiFvpj/JDgKKd7ikh4
  sudo: ALL=(ALL) NOPASSWD:ALL

ssh_pwauth: true

timezone: Etc/UTC

runcmd:
- localectl set-x11-keymap "us" pc105
- setupcon -k --force || true
- astroberry-init
EOF

#if [ -e /boot/firmware/network-config ]; then cp /boot/firmware/network-config /boot/firmware/network-config.default; fi
#cat <<EOF > /boot/firmware/network-config
#network:
#  version: 2
#
#  ethernets:
#    eth0:
#      dhcp4: true
#      optional: true
#
#  wifis:
#    wlan0:
#      dhcp4: no
#      addresses: [10.42.0.1/24]
#      access-points:
#        "astroberry":
#          password: "astroberry"
#          mode: ap
#      regulatory-domain: GB
#EOF

# Set wireless regulatory domain
if [ -e /boot/firmware/cmdline.txt ] && [ -z "$(grep cfg80211.ieee80211_regdom /boot/firmware/cmdline.txt)" ]; then
	sed -i -e "s/\s*cfg80211.ieee80211_regdom=\S*//" -e "s/\(.*\)/\1 cfg80211.ieee80211_regdom=GB/" /boot/firmware/cmdline.txt
fi

# Set reverse proxy / caddy
if [ -e /etc/caddy/Caddyfile ]; then cp /etc/caddy/Caddyfile /etc/caddy/Caddyfile.default; fi
cat <<EOF > /etc/caddy/Caddyfile
{
	auto_https disable_redirects
}

http://astroberry.local {
        handle /api* {
            reverse_proxy localhost:8624
        }
        handle /websockify* {
            reverse_proxy localhost:8070
        }
        handle /* {
            reverse_proxy localhost:8080
        }
}

https://astroberry.local {
        handle /api* {
            reverse_proxy localhost:8624
        }
        handle /websockify* {
            reverse_proxy localhost:8070
        }
        handle /* {
            reverse_proxy localhost:8080
        }
}
EOF

# Set filesharing / samba
if [ -e /etc/samba/smb.conf ]; then cp /etc/samba/smb.conf /etc/samba/smb.conf.default; fi
cat <<EOF > /etc/samba/smb.conf
[global]
   workgroup = WORKGROUP
   log file = /var/log/samba/log.%m
   max log size = 1000
   logging = file
   panic action = /usr/share/samba/panic-action %d
   server role = standalone server
   obey pam restrictions = yes
   unix password sync = yes
   passwd program = /usr/bin/passwd %u
   passwd chat = *Enter\snew\s*\spassword:* %n\n *Retype\snew\s*\spassword:* %n\n *password\supdated\ssuccessfully* .
   pam password change = yes
   map to guest = bad user
   usershare allow guests = no
   multicast dns register = no

[Shared Files]
   comment = Shared Files
   path = /home/astroberry/Shared
   browseable = yes
   writeable = yes
   read only = no
   create mask = 0644
   directory mask = 0755
   valid users = astroberry
EOF

# Install astroberry-manager
if [ ! -e /opt/astroberry-manager ]; then python -m venv /opt/astroberry-manager; fi
/opt/astroberry-manager/bin/pip install --upgrade astroberry_manager@git+https://github.com/astroberry-official/astroberry-manager

# Disable graphical login
# update-rc.d lightdm disable

# Enable TigerVNC
# ln -s /etc/systemd/system/vncserver.service /etc/systemd/system/multi-user.target.wants/

# Enable INDI Web Manager
# ln -s /etc/systemd/system/indiwebmanager.service /etc/systemd/system/multi-user.target.wants/

# Enable Websockify
# ln -s /etc/systemd/system/websockify.service /etc/systemd/system/multi-user.target.wants/

# Enable Astroberry Manager
# ln -s /etc/systemd/system/astroberry-manager.service /etc/systemd/system/multi-user.target.wants/

fi
