# debian/prerm
#!/bin/bash
# prerm script for astroberry-os-sysmod

if [ "$1" == "remove" -o "$1" == "purge" ]; then

        # Disable services
        systemctl stop astroberry-manager.service && systemctl disable astroberry-manager.service
	if [ -e /opt/astroberry-manager ]; then rm -rf /opt/astroberry-manager; fi
        systemctl stop indiwebmanager.service && systemctl disable indiwebmanager.service
        systemctl stop vncserver.service && systemctl disable vncserver.service
        systemctl stop websockify.service && systemctl disable websockify.service

fi
